# agents/zu6ydrahph1tt6wxw7jpbnmzw1gk5hu/com.advectus.zu6ydrahph1tt6wxw7jpbnmzw1gk5hu.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
   <key>Label</key>
   <string>com.advectus.zu6ydrahph1tt6wxw7jpbnmzw1gk5hu</string>
   <key>Program</key>
   <string>/Users/edmricha/Development/MyLibrary/osx-launch-agents/agents/zu6ydrahph1tt6wxw7jpbnmzw1gk5hu/exec.sh</string>
   <key>RunAtLoad</key>
   <true/>
</dict>
</plist>

